# src/systemd/systemd-backlight@.service.d/zfs-rootfs.conf
#ZFS related dependencies
[Unit]
After=zfs-mount.service
