# src/systemd/zfs-mount.service.d/root-fs-dependencies.conf
#ZFS related dependencies
[Unit]
Before=systemd-rfkill.service systemd-random-seed.service
